# native/pan-bridge-vps/pan-bridge.service
[Unit]
Description=EI Solutions — PAN Bridge (HMAC proxy → mallikacyberzone)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=panbridge
WorkingDirectory=/opt/ei-pan-bridge
EnvironmentFile=/opt/ei-pan-bridge/.env
ExecStart=/usr/bin/node /opt/ei-pan-bridge/server.js
Restart=always
RestartSec=5
StandardOutput=append:/var/log/ei-pan-bridge/out.log
StandardError=append:/var/log/ei-pan-bridge/err.log

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/ei-pan-bridge

[Install]
WantedBy=multi-user.target
