# native/bbps-bridge-vps/bbps-bridge.service
[Unit]
Description=EI SOLUTIONS — Bharat Connect (BBPS) HMAC Proxy Bridge
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/bbps-bridge
EnvironmentFile=/opt/bbps-bridge/.env
ExecStart=/usr/bin/node /opt/bbps-bridge/server.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
